# envy13t-check/check.service
[Unit]
Description=Check services on HP Envy 13t

[Service]
ExecStart=/usr/bin/envy13t-check

[Install]
WantedBy=default.target
